# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod --config /etc/mongod.conf

[program:redis]
command=/opt/local/redis-stable/src/redis-server /etc/redis/6379.conf